# ha-pluggie/rootfs/etc/services.d/admin_api/run
#!/usr/bin/with-contenv bashio

export LOG_LEVEL=$(bashio::config 'log_level' 'info')

/usr/local/bin/admin_api.py
